# PROCEDURES/SWM/tools/2.0/mi_cnz/esm/virtual/vmware/apg_dump/apg-cluster-update.service
[Unit]
Description=Sync MAC addresses between control nodes
After=lde-failoverd.service lde-nfs-client.service
Before=lde-rpm.service

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/apg_cluster_update.sh start
RemainAfterExit=yes

[Install]
WantedBy=lde-nodetype@control.target
