# PROCEDURES/SWM/tools/2.0/mi_cnz/esm/virtual/vmware/apg_dump/apos-early-system-config.service
[Unit]
Description=APOS early system configuration
Before=lde-control.service lde-network.service lde-cloud-init.service

[Service]
Type=oneshot
ExecStart=/usr/lib/systemd/scripts/apos-system-conf.sh early
RemainAfterExit=yes

[Install]
RequiredBy=lde-control.service lde-network.service lde-cloud-init.service
